# src/tetex-1.0.5-1cl/tetex.cron
#!/bin/bash
# Remove and TeX fonts not used in 10 days
/usr/sbin/tmpwatch 240 /var/lib/texmf
exit 0


